# files/system/usr/lib/systemd/system/sddm.service.d/10-ensure-sysusers.conf
[Service]
# Ensure the sddm system user exists before the greeter starts.
# systemd-sysusers.service has ConditionNeedsUpdate=/etc and is skipped on
# subsequent boots of the same deployment, leaving the sddm user absent and
# causing a black screen. Running sysusers unconditionally here fixes that.
ExecStartPre=/usr/bin/systemd-sysusers /usr/lib/sysusers.d/sddm.conf
